# systemd/unbound-ui-doh-proxy.service
[Unit]
Description=Cloudflared DoH Proxy for Pusula
Documentation=https://developers.cloudflare.com/cloudflare-one/connections/connect-apps/
After=network.target

[Service]
Type=simple
User=cloudflared
ExecStart=/usr/local/bin/cloudflared proxy-dns --port 5053 --upstream https://cloudflare-dns.com/dns-query

# Restart policy
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cloudflared-doh

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
